# supervisord.conf
[supervisord]
nodaemon=true

[program:mariadb]
command=/usr/bin/mysqld_safe --basedir=/usr & sleep 3s
